# debian/etc/systemd/system/getty@tty1.service.d/dietpi-autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --noclear --noissue -a dietpi -J %I $TERM